# ansible/roles/deploy/templates/utg-bot.service.j2
[Unit]
Description=UTG Discord Bot

[Service]
Type=simple
ExecStart=/usr/bin/utg-bot
Restart=always

[Install]
WantedBy=multi-user.target

